# v3_1/dockerfiles/build/configs/supervisord/orocommerce.conf
[supervisord]
nodaemon = true
loglevel = INFO
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0700
username = dummy
password = dummy

[supervisorctl]
username = dummy
password = dummy

[program:php-fpm]
command=/usr/sbin/php-fpm7.1 -F
stopasgroup=true
stopsignal=TERM
autostart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:cron]
command=/usr/sbin/cron -f
stopsignal=TERM
stopasgroup=true
autostart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true


[program:oro_web_socket]
command=php ./app/console clank:server --env=prod
numprocs=1
autostart=true
autorestart=true
directory=/home/orocommerce
user=orocommerce
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:oro_message_consumer]
command=php ./app/console oro:message-queue:consume --env=prod
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
directory=/home/orocommerce
user=orocommerce
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
